# product/images/packetbeat/supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslogd]
command=/usr/sbin/rsyslogd

[program:packetbeat]
command=/usr/bin/packetbeat -c /etc/packetbeat/packetbeat.yml