# addons/packetbeat/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=rsyslog

[program:packetbeat]
command=/etc/init.d/packetbeat start